# systemd/quail-purge.timer
[Unit]
Description=Quail daily retention purge

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
